# resources/platform/termux/data/data/com.termux/files/usr/var/service/plug/run
#!/data/data/com.termux/files/usr/bin/sh
cd /data/data/com.termux/files/usr/opt/plug
exec node service/main.js 2>&1